# build/etc/supervisor/conf.d/mailcatcher.conf
[program:mailcatcher]
command=/usr/local/bin/mailcatcher --ip 0.0.0.0 --smtp-port 25 -f
autostart=true
autorestart=true