# build/etc/supervisor/conf.d/mailcatcher.conf
[program:mailcatcher]
command=/usr/local/bin/mailcatcher -f
autostart=true
autorestart=true