# src/build/modules/mailcatcher/files/etc/supervisor/conf.d/mailcatcher.conf
[program:mailcatcher]
command = /usr/local/bin/mailcatcher --ip 0.0.0.0 --http-port 1080 --smtp-port 1025 -f
autostart = true
autorestart = true
